# nginx.conf
[program:nginx]
command = /sbin/nginx -g "daemon off;" 
user = root  
autostart = true
stderr_logfile = /dev/stdout
stdout_logfile = /dev/stdout
