# php.conf
[program:php]
command = /sbin/php-fpm  -F
user = root  
autostart = true 
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout 
